# cron/crontab
*/1 * * * * /usr/local/bin/python /app/start_memtester.py > /proc/1/fd/1 2>/proc/1/fd/2
*/2 * * * * /usr/local/bin/python /app/fetch_probe_stats.py > /proc/1/fd/1 2>/proc/1/fd/2
